# config/supervisord/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[supervisord]
nodaemon=true

;[include]
;files = /home/config/supervisord/*.conf

; program-sample in /home/config/supervisord/*.conf
; ** please check logfile is existed **
; --------------------------------------------------------
;[program:laravel-worker]  
;process_name=%(program_name)s_%(process_num)02d
;user=root
;command=php /home/website/kr8851queue_web/artisan queue:work --queue=csvsite --tries=5 --memory=512 --sleep=5
;autostart=true
;autorestart=true
;stderr_logfile=/home/logs/supervisord/laravel-worker/error.log
;stdout_logfile=/home/logs/supervisord/laravel-worker/nohup.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
stdout_logfile=/home/logs/nginx/daemon_out.log
stderr_logfile=/home/logs/nginx/daemon_err.log

[program:php-fpm]
command=/usr/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
